# maint.conf
[program:maint]
command=su - www-data -c "/usr/bin/php /var/www/html/timetrex/maint/cron.php;sleep 60;"
stdout_events_enabled=true
stderr_events_enabled=true
startretries=0
autorestart=true
