# postgres.conf
[program:postgresql]
command=/usr/lib/postgresql/9.5/bin/postgres -c config_file=/etc/postgresql/9.5/main/postgresql.conf
user=postgres
autorestart=true
